# server-config/crontab/crontab.ct
45 * * * * php /var/www/html/trickoreat/backendtoe/cron-jobs/clean-reset-tokens/run.php